# DockerImage/rpm-repository/spacewalk/CentOS6/supervisord.conf
[supervisord]
nodaemon=true

[program:spacewalk]
command=/usr/sbin/spacewalk-service start

