# DockerImage/load-balancer/centos6/supervisord.d/piranha-gui.ini
[program:piranha-gui]
command=service piranha-gui start
