# DockerImage/load-balancer/centos6/supervisord.d/pulse.ini
[program:pulse]
command=service pulse start
