# DockerImage/rpm-repository/CentOS6/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/spacewalk-service start

